# etc/cron.d/cfbackup
0 4 * * *	root	/opt/Colfusion/ColFusion/assets/scripts/backup.sh
