# etc/cron.d/cfopenrefine
@reboot	root	sleep 60 && /etc/init.d/cfopenrefine start
